# dev/feather.rc
#!/bin/sh

# PROVIDE: feather
# REQUIRE: LOGIN DAEMON NETWORKING
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable feather:
#
#   feather_enable="YES"
#
# Optional rc.conf variables:
#
#   feather_dir="/usr/local/feather"       # Installation directory
#   feather_user="root"                    # User to run as (root needed for port <1024)
#   feather_group="wheel"                  # Group to run as
#   feather_config="/usr/local/etc/feather" # Config directory
#   feather_log="/var/log/feather"         # Log directory

. /etc/rc.subr

name="feather"
rcvar="${name}_enable"

load_rc_config $name

: ${feather_enable:="NO"}
: ${feather_dir:="/usr/local/feather"}
: ${feather_user:="root"}
: ${feather_group:="wheel"}
: ${feather_config:="/usr/local/etc/feather"}
: ${feather_log:="/var/log/feather"}

pidfile="/var/run/${name}.pid"
feather_bin="${feather_dir}/bin/feather"

export LANG="en_US.UTF-8"
export RELEASE_NODE="feather@localhost"
export FEATHER_CONFIG_FOLDER="${feather_config}"
export RELEASE_TMP="${feather_dir}/tmp"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

feather_start()
{
    if feather_running; then
        echo "${name} is already running."
        return 0
    fi

    echo "Starting ${name}."
    install -d -o "${feather_user}" -g "${feather_group}" "${feather_log}"
    install -d -o "${feather_user}" -g "${feather_group}" "${feather_dir}/tmp"

    # Use daemon(8) to fork the BEAM into the background, write its PID to
    # ${pidfile}, redirect stdio to console.log, and (if feather_user != root)
    # drop privileges. We invoke `feather start` (not `feather daemon`) so the
    # BEAM stays attached to daemon(8) as a normal child — `feather stop`
    # signals it cleanly and daemon(8) removes the pidfile on exit.
    /usr/sbin/daemon \
        -p "${pidfile}" \
        -u "${feather_user}" \
        -o "${feather_log}/console.log" \
        env HOME="${feather_dir}" \
            LANG="${LANG}" \
            RELEASE_NODE="${RELEASE_NODE}" \
            FEATHER_CONFIG_FOLDER="${FEATHER_CONFIG_FOLDER}" \
            RELEASE_TMP="${RELEASE_TMP}" \
            "${feather_bin}" start

    # Wait for the node to register with epmd. BEAM cold-start can take a few
    # seconds, so poll instead of using a fixed sleep.
    _i=0
    while [ $_i -lt 30 ]; do
        if feather_running; then
            echo "${name} started."
            return 0
        fi
        sleep 1
        _i=$((_i + 1))
    done

    echo "${name} failed to start. Last lines from console log:"
    [ -f "${feather_log}/console.log" ] && tail -60 "${feather_log}/console.log"
    return 1
}

feather_stop()
{
    if ! feather_running; then
        echo "${name} is not running."
        rm -f "${pidfile}"
        return 0
    fi

    echo "Stopping ${name}."
    # Graceful shutdown via the release's own RPC. daemon(8) sees the BEAM
    # exit and tears down on its own.
    HOME="${feather_dir}" "${feather_bin}" stop 2>/dev/null || true

    _i=0
    while [ $_i -lt 20 ]; do
        if ! feather_running; then
            rm -f "${pidfile}"
            echo "${name} stopped."
            return 0
        fi
        sleep 1
        _i=$((_i + 1))
    done

    # RPC didn't take it down — fall back to SIGTERM on the tracked pid.
    echo "${name} did not stop within 20s; sending TERM to pid $(cat "${pidfile}" 2>/dev/null)."
    [ -f "${pidfile}" ] && kill -TERM "$(cat "${pidfile}")" 2>/dev/null
    sleep 2
    rm -f "${pidfile}"
}

feather_status()
{
    if feather_running; then
        _pid=$(feather_pid)
        echo "${name} is running (pid ${_pid})."
    else
        echo "${name} is not running."
        return 1
    fi
}

feather_running()
{
    HOME="${feather_dir}" "${feather_bin}" pid > /dev/null 2>&1
}

feather_pid()
{
    HOME="${feather_dir}" "${feather_bin}" pid 2>/dev/null
}

run_rc_command "$1"
