# dev/feather.rc
#!/bin/sh

# PROVIDE: feather
# REQUIRE: LOGIN DAEMON NETWORKING
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable feather:
#
#   feather_enable="YES"
#
# Optional rc.conf variables:
#
#   feather_dir="/usr/local/feather"       # Installation directory
#   feather_user="root"                    # User to run as (root needed for port <1024)
#   feather_group="wheel"                  # Group to run as
#   feather_config="/usr/local/etc/feather" # Config directory
#   feather_log="/var/log/feather"         # Log directory

. /etc/rc.subr

name="feather"
rcvar="${name}_enable"

load_rc_config $name

: ${feather_enable:="NO"}
: ${feather_dir:="/usr/local/feather"}
: ${feather_user:="root"}
: ${feather_group:="wheel"}
: ${feather_config:="/usr/local/etc/feather"}
: ${feather_log:="/var/log/feather"}

pidfile="/var/run/${name}.pid"
feather_bin="${feather_dir}/bin/feather"

export LANG="en_US.UTF-8"
export RELEASE_NODE="feather@localhost"
export FEATHER_CONFIG_FOLDER="${feather_config}"
export RELEASE_TMP="${feather_dir}/tmp"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

feather_start()
{
    if feather_running; then
        echo "${name} is already running."
        return 0
    fi

    echo "Starting ${name}."
    install -d -o "${feather_user}" -g "${feather_group}" "${feather_log}"
    install -d -o "${feather_user}" -g "${feather_group}" "${feather_dir}/tmp"

    # feather is a system user with no real $HOME; point it at the install dir
    # (which it owns) so the release boot script and run_erl have somewhere to
    # write. Without this, su -m preserves root's HOME and the BEAM dies before
    # the application logger comes up.
    HOME="${feather_dir}" su -m "${feather_user}" -c "${feather_bin} daemon"

    # Wait for the node to register with epmd. BEAM cold-start can take a few
    # seconds, so poll instead of using a fixed sleep.
    _i=0
    while [ $_i -lt 30 ]; do
        if feather_running; then
            feather_pid > "${pidfile}"
            echo "${name} started."
            return 0
        fi
        sleep 1
        _i=$((_i + 1))
    done

    echo "${name} failed to start. Last lines from boot log:"
    for _f in "${feather_dir}/tmp/log/erlang.log."*; do
        [ -f "$_f" ] || continue
        echo "---- $_f ----"
        tail -40 "$_f" 2>/dev/null
    done
    return 1
}

feather_stop()
{
    if ! feather_running; then
        echo "${name} is not running."
        return 0
    fi

    echo "Stopping ${name}."
    HOME="${feather_dir}" su -m "${feather_user}" -c "${feather_bin} stop"
    rm -f "${pidfile}"
    echo "${name} stopped."
}

feather_status()
{
    if feather_running; then
        _pid=$(feather_pid)
        echo "${name} is running (pid ${_pid})."
    else
        echo "${name} is not running."
        return 1
    fi
}

feather_running()
{
    HOME="${feather_dir}" su -m "${feather_user}" -c "${feather_bin} pid" > /dev/null 2>&1
}

feather_pid()
{
    HOME="${feather_dir}" su -m "${feather_user}" -c "${feather_bin} pid" 2>/dev/null
}

run_rc_command "$1"
